# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

; Control socket so `supervisorctl status` works (e.g. in a Render shell) to
; confirm php-fpm, nginx and the queue worker are RUNNING.
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Queued mail + notifications. Same container is fine for a single-instance
; deploy; scale this out to a dedicated Render Background Worker if needed.
[program:queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --no-interaction
autostart=true
autorestart=true
user=www-data
; Run as www-data but with an accessible HOME. Otherwise HOME stays /root
; (from the root supervisord) and libpq tries to read the optional client cert
; at /root/.postgresql/postgresql.crt — which www-data can't access — failing
; every Neon SSL connection with "Permission denied".
directory=/var/www/html
environment=HOME="/tmp",USER="www-data"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Laravel Reverb WebSocket server (chat / notifications / call signaling).
; Binds to localhost; nginx proxies wss traffic to it on the public port so the
; whole app runs on Render's single $PORT. Same HOME fix as the queue worker.
[program:reverb]
command=php /var/www/html/artisan reverb:start --host=127.0.0.1 --port=8080 --no-interaction
autostart=true
autorestart=true
user=www-data
directory=/var/www/html
environment=HOME="/tmp",USER="www-data"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
